# package/harvester-os/files/etc/systemd/system/rancherd.service.d/override.conf
[Unit]
After=time-sync.target
Wants=time-sync.target
