# package/harvester-os/files/etc/systemd/system/kdump.service.d/override.conf
[Service]
ExecStartPre=/usr/sbin/gen-kdump-initrd
