# package/harvester-os/files/etc/systemd/system/rancher-system-agent.service.d/env.conf
[Service]
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/opt/rke2/bin"
